# config/autostart/per-de-config.desktop
[Desktop Entry]
Type=Application
Exec=sh -c "~/.bin/per-de-config.sh"
