# config/autostart/per-de-config.desktop
[Desktop Entry]
Name=Per DE Config Files
Type=Application
Exec=sh -c "~/.bin/per-de-config.sh"
